# configuration/home/.config/systemd/user/connect-bt.service
[Unit]
Description=Bluetooth devices etc.

[Service]
ExecStart=/home/dreuter/.scripts/connect_bt.sh

[Install]
WantedBy=default.target
