# configuration/home/.config/systemd/user/dev-startup.service
[Unit]
Description=Notebooks etc.

[Service]
ExecStart=/home/dreuter/.scripts/dev-startup.sh

[Install]
WantedBy=default.target
