# configuration/home/.config/systemd/user/hwmon-link.service
[Unit]
Description=Search and link hwmon cpu temp

[Service]
ExecStart=/home/dreuter/.scripts/hwmon-link.sh

[Install]
WantedBy=default.target
